# findmore/galp/Bash/crontab
0 2 * * * find /* -mtime +182 -exec rm {} +